# assets/systemd-units/adjust_rgb_brightness_by_time.timer
[Unit]
Description=Adjust the brightness of RGB devices at 18:00

[Timer]
OnCalendar=18:00

[Install]
WantedBy=timers.target

